# com.thyself.weekly-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thyself.weekly-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jfru/thyself/sync/run_sync.sh</string>
    </array>

    <!-- Run every Sunday at 3:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jfru/Library/Application Support/Thyself/logs/sync-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jfru/Library/Application Support/Thyself/logs/sync-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/jfru/thyself</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jfru/.pyenv/shims:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
